# docker/rootfs/etc/services.d/alice-transmission-dialogs/run
#!/usr/bin/with-contenv sh
alice_transmission
